# scripts/claude-epaper.service
[Unit]
Description=Claude e-paper — app web + API (Node)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Adapte User et les chemins à ton installation.
User=pi
WorkingDirectory=/home/pi/claude-epaper
Environment=NODE_ENV=production
Environment=PORT=8787
ExecStart=/usr/bin/node server/dist/index.js
Restart=always
RestartSec=3
# KillMode=process : au restart, systemd ne tue que le process principal — le
# script de mise à jour (lancé depuis le dashboard) survit et peut finir son job.
KillMode=process

[Install]
WantedBy=multi-user.target
