# scripts/epaper-push.service
[Unit]
Description=Claude e-paper — pousse /api/render.png sur la dalle
After=network-online.target claude-epaper.service
Wants=network-online.target

[Service]
Type=simple
# Adapte User et les chemins à ton installation.
User=pi
WorkingDirectory=/home/pi/claude-epaper
# Module waveshare_epd : chemin de la lib officielle clonée (install native).
Environment=PYTHONPATH=/home/pi/e-Paper/RaspberryPi_JetsonNano/python/lib
# Variables : soit ici en clair, soit via EnvironmentFile (décommenter).
Environment=RENDER_URL=http://localhost:8787/api/render.png?palette=bw
Environment=EPD_MODEL=epd2in13_V4
Environment=POLL_SECONDS=5
Environment=FULL_REFRESH_EVERY=20
Environment=FULL_REFRESH_SECONDS=1800
#EnvironmentFile=/home/pi/claude-epaper/scripts/epaper-push.env
ExecStart=/usr/bin/python3 /home/pi/claude-epaper/scripts/epaper_push.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
