# scripts/epaper-push.service
[Unit]
Description=Claude e-paper — pousse /api/render.png sur la dalle
After=network-online.target claude-epaper.service
Wants=network-online.target

[Service]
Type=simple
# Adapte User et les chemins à ton installation.
User=pi
WorkingDirectory=/home/pi/claude-epaper
# Module waveshare_epd : chemin de la lib officielle clonée (install native).
Environment=PYTHONPATH=/home/pi/e-Paper/RaspberryPi_JetsonNano/python/lib
# Seulement les variables spécifiques à la machine ; les réglages de comportement
# (POLL_SECONDS, FULL_REFRESH_*, EPAPER_THRESHOLD…) suivent les défauts du script
# → une mise à jour du code les fait évoluer sans refaire `make services`.
Environment=RENDER_URL=http://localhost:8787/api/render.png
Environment=EPD_MODEL=epd2in13_V4
#EnvironmentFile=/home/pi/claude-epaper/scripts/epaper-push.env
ExecStart=/usr/bin/python3 /home/pi/claude-epaper/scripts/epaper_push.py
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
